# zabbix-apcups.cron
# Zabbix APC UPS Check:
*/1 * * * * zabbix /etc/zabbix/scripts/zabbix_apcups.sh 1>/dev/null 2>/dev/null
